# homepage-contact-count.service
[Unit]
Description=Update homepage contact unread count

[Service]
Type=oneshot
ExecStart=/bin/bash /home/kdy987/law-crawler/ui/homepage-contact-count.sh
User=kdy987