# law-crawler-ui.service
[Unit]
Description=Law Crawler UI Service
After=network.target

[Service]
Type=simple
User=kdy987
Group=kdy987
WorkingDirectory=/home/kdy987/law-crawler/ui
# 가상환경의 python을 직접 사용하여 uvicorn을 모듈로 실행 (-m)
ExecStart=/home/kdy987/law-crawler/ui/.venv/bin/python -m uvicorn app.backend.main:app --host 0.0.0.0 --port 8004

# 환경 변수 확인: T3600은 CachyOS이므로 필요한 경우 수정하세요
Environment="LAW_CRAWLER_MODE=fedora"
Environment="PYTHONPATH=/home/kdy987/law-crawler/ui"

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target